# docker/devbox-autoconf/etc/service/devbox-autoconfig/run
#!/bin/bash
devbox-autoconf /etc/devbox-autoconf/config.ini \
  2> >(log-prefix \! autoconf) \
  1> >(log-prefix . autoconf)
