# docker/php/etc/service/php7.4-fpm/run
#!/bin/bash
php-fpm7.4 \
  --nodaemonize \
  --fpm-config /etc/php/7.4/fpm/php-fpm.conf \
  2> >(log-prefix \! php-fpm7.4) \
  1> >(log-prefix . php-fpm7.4)
