# docker/files.apt-post/etc/service/ssh-agent/run
#!/bin/bash
chpst -u www-data ssh-agent -a /run/ssh-agent/.sock -D \
  2> >(log-prefix \! ssh-agent) \
  1> /dev/null
