# docker/files.apt-post/etc/service/ssh-agent/run
#!/bin/sh
exec chpst -u www-data ssh-agent -a /run/ssh-agent/.sock -D
