# docker/mariadb/post-install/etc/service/mariadb/run
#!/bin/bash
mysqld \
  --init-file=/mariadb_init.sql \
  --silent-startup \
  2> >(log-prefix \! mariadb) \
  1> >(log-prefix . mariadb)
