# docker/base/etc/service/ssh-agent/run
#!/bin/bash
install -m 750 -d /run/ssh-agent -o devbox -g devbox
chpst -u devbox ssh-agent -a /run/ssh-agent/.sock -D \
  2> >(log-prefix \! ssh-agent) \
  1> /dev/null
