# docker/base/etc/service/ssh-agent/run
#!/bin/bash
ssh-agent -a /run/ssh-agent-devbox.socket -D \
  2> >(log-prefix \! ssh-agent) \
  1> /dev/null
