# docker/devbox-py/etc/service/devbox-autoconfig/run
#!/bin/bash
devbox-autoconf \
  2> >(log-prefix \! autoconf) \
  1> >(log-prefix . autoconf)
