# docker/mailcatcher/post-install/etc/service/mailcatcher/run
#!/bin/bash
mailcatcher \
  --foreground \
  --http-ip 127.0.0.1 \
  --no-quit \
  &> /var/log/mailcatcher.log
