# docker/files.apt-post/etc/service/mailcatcher/run
#!/bin/bash
mailcatcher \
  --foreground \
  --http-ip 0.0.0.0 \
  --no-quit \
  &> /var/log/mailcatcher.log
