# docker/files.apt-post/etc/service/apache2/run
#!/bin/bash
apachectl -D FOREGROUND \
  2> >(log-prefix \! apache2) \
  1> >(log-prefix . apache2)
