# docker/files.apt-post/etc/service/devbox-autoconfig/run
#!/bin/bash
devbox-autoconfig.py /etc/devbox-autoconfig.ini \
  2> >(log-prefix \! autoconfig) \
  1> >(log-prefix . autoconfig)
