# docker/files.apt-post/etc/service/devbox-autoconfig/run
#!/bin/sh
exec /usr/local/sbin/devbox-autoconfig.py /etc/devbox-autoconfig.ini
