# docker/files.apt-post/etc/service/mariadb/run
#!/bin/sh
exec mysqld --user=mysql --init-file=/mariadb_init.sql
